# Server/etc/systemd/noise-processor.service
[Unit]
Description=Noise Data Processor Service
After=network.target mosquitto.service
Requires=mosquitto.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/noise-backend
ExecStart=/usr/bin/python3 /home/pi/noise-backend/mqtt_processor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/noise-backend
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
